# rr.conf
[program:rr]
command=python3 /home/curator2/rr/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/rr/test.err.log
stdout_logfile=/home/curator2/rr/test.out.log